# deploy/systemd/mqcmdb-pipeline.timer
[Unit]
Description=Run MQ CMDB Pipeline daily

[Timer]
# Run daily at 6:00 AM
OnCalendar=*-*-* 06:00:00

# Add some randomness to avoid thundering herd
RandomizedDelaySec=300

# Run immediately if we missed the last scheduled time
Persistent=true

[Install]
WantedBy=timers.target
